# socketio/supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log		; main log file; default $CWD/supervisord.log
nodaemon=true                  			; start in foreground if true; default false
pidfile=/run/supervisord.pid    		; supervisord pidfile; default supervisord.pid
user=root                    			; setuid to this UNIX account at startup; recommended if root

; The [rpcinterface:supervisor] section must
; remain in the configuration for the standard
; setup of supervisor to work properly.
; http://supervisord.org/configuration.html
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock	; use a unix:// URL for a unix socket

[program:pong-socketio]
command=/bin/bash /app/pong_socketio/entrypoint.sh
redirect_stderr=true					; If true, cause the process’ stderr output to be sent back to supervisord on its stdout file
autorestart=false						; Specifies if supervisord should automatically restart a process if it exits when it is in the RUNNING state. 
startretries=2							

[program:chat-socketio]
command=/bin/bash /app/chat_socketio/entrypoint.sh
redirect_stderr=true					; If true, cause the process’ stderr output to be sent back to supervisord on its stdout file
autorestart=false						; Specifies if supervisord should automatically restart a process if it exits when it is in the RUNNING state. 
startretries=2		